# init.d/checkroot
#!/sbin/runscript
# Copyright 1999-2002 Gentoo Technologies, Inc.
# Distributed under the terms of the GNU General Public License, v2 or later
# $Header$


depend() {
	before *
}

start() {
	ebegin "Remounting root filesystem read-only (if necessary)"
	mount / -o remount,ro
	eend $?
	ebegin "Checking root filesystem"
	fsck -a /
	if [ $? -eq 0 ]
	then
		eend 0
	elif [ $? -eq 1 ]
	then
		ewend 1 "Filesystem repaired"
	else
		eend 2 "Filesystem couldn't be fixed :("
		/sbin/sulogin $CONSOLE
		reboot -f
	fi

	ebegin "Remounting root filesystem read/write"
	mount / -o remount,rw
	eend $?

	if [ "$BOOT" = "yes" ]
	then
		# create /etc/mtab (done by forcing current mounted filesystems)
		# NB:  this needs to be updated if our /sbin/rc mounts change ...
		> /etc/mtab
		mount -f /
		mount -f /proc

		# should /dev be mounted ?
		get_bootparam "nodevfs"
		if [ $? -ne 0 ] && [ -e /dev/.devfsd ]
		then
			mount -f -t devfs none /dev
		fi
		
		mount -f -t tmpfs tmpfs /dev/shm -o rw,mode=0644,size=4m
	fi
	
	return
}

