# init.d/checkfs
#!/sbin/runscript
# Copyright 1999-2002 Gentoo Technologies, Inc.
# Distributed under the terms of the GNU General Public License, v2 or later
# $Header$


depend() {
	need checkroot
}

start() {
	# LVM support for /usr, /home, /opt ....
	# This should be done *before* checking local
	# volumes, or they never get checked.
	if [ -x /sbin/vgscan -a -d /proc/lvm ]
	then
		ebegin "Setting up the Logical Volume Manager"
		#still echo stderr for debugging
		/sbin/vgscan >/dev/null 
		if [ "$?" -eq 0 ] && [ -x /sbin/vgchange -a -f /etc/lvmtab ]
		then
			/sbin/vgchange -a y >/dev/null
		fi
		eend $? "Failed to setup the LVM"
	fi

	# Start software raid
	if [ -x /sbin/raidstart ] && [ -f /etc/raidtab ] && [ -f /proc/mdstat ]
	then
		if [ `grep "persistent-superblock" /etc/raidtab | awk '{print $2}'` == 0 ]
		then
			ebegin "Starting software RAID"
			# still echo stderr for debugging
			/sbin/raidstart --all >/dev/null
			eend $? "Failed to start software RAID"
		fi
	fi

	if [ -f /fastboot ]
	then
		rm -f /fastboot
	else
		ebegin "Checking all filesystems"
		if [ -f /forcefsck ]
		then
			ewarn "A full fsck has been forced"
			fsck -C -R -A -a -f
			rm -f /forcefsck
		else
			fsck -R -A -a
		fi
		if [ "$?" -eq 0 ]
		then
			eend 0
		elif [ "$?" -eq 1 ]
		then
			ewend 1 "Filesystem errors corrected."
			#everything should be ok, so return a pass
			return 0
		else
			eend 2 "Fsck could not correct all errors, manual repair needed"
			/sbin/sulogin ${CONSOLE}
		fi
	fi
}


# vim:ts=4
