# init.d/checkfs
#!/sbin/runscript

depend() {
	need checkroot
}

start() {
	ebegin "Checking all filesystems"
	fsck -R -A -a
	if [ $? -eq 0 ]
	then
		eend 0
	elif [ $? -eq 1 ]
	then
		eend 1 "Filesystem errors corrected."
	else 
		eend 2 "Fsck could not correct all errors, manual repair needed"
		/sbin/sulogin $CONSOLE
	fi
	return
}
