# init.d/consolefont
#!/sbin/runscript
# Copyright 1999-2002 Gentoo Technologies, Inc.
# Distributed under the terms of the GNU General Public License, v2 or later
# $Header$


depend() {
	need localmount
}

start() {
	if [ -n "$CONSOLEFONT" ] ; then
		if [ -x /usr/bin/setfont ] ; then
			local param x
			if [ -n "$CONSOLETRANSLATION" ] ; then
				param="-m ${CONSOLETRANSALTION}"
			fi
			for x in 1 2 3 4 5 6
			do
				if [ -e /dev/.devfsd ]
				then
					/usr/bin/setfont ${CONSOLEFONT} ${param} \
						--tty=/dev/vc/${x}
				else
					/usr/bin/setfont ${CONSOLEFONT} ${param} \
						--tty=/dev/tty${x}
				fi
			done
		fi
	fi
	return
}
