# init.d/consolefont
#!/sbin/runscript

depend() {
	need localmount
}

start() {
	source /etc/rc.conf
	if [ -n "$CONSOLEFONT" ]
	then
		if [ -e /usr/bin/consolechars ]
		then
			/usr/bin/consolechars --font=${CONSOLEFONT}
		fi
	fi
	return
}

