# init.d/devfsd
#!/sbin/runscript
SERVICE="devfsd"
opts="start stop"

start() {
	need pretty
	if [ ! -e /dev/.devfsd ]
	then
		return
	fi
	#the "sleep" is to allow devfsd to start up and make compatibility
	#devices ready to the fsck routines
	ebegin "Starting devfsd"; /sbin/devfsd /dev; sleep 2; eend $?
}

stop() {
	ebegin "Stopping devfsd"; killall devfsd; eend $?
}
