# init.d/hostname
#!/sbin/runscript
# Copyright 1999-2002 Gentoo Technologies, Inc.
# Distributed under the terms of the GNU General Public License, v2 or later
# $Header$


depend() {
	need localmount
}

checkconfig() {
	if [ -z "`cat /etc/hostname`" ] ; then
		eerror "You need to set /etc/hostname to a valid hostname"
		return 1
	fi
}

start() {
	if [ ! checkconfig ] ; then
		myhost=localhost
	else
		myhost=`cat /etc/hostname` 1>&2
	fi
	ebegin "Setting hostname to $myhost"
	/bin/hostname ${myhost}
	eend $?
	return
}
