# init.d/keymaps
#!/sbin/runscript
	
depend() {
	need localmount
}

start() {
	ebegin "Loading key mappings"
	/usr/bin/loadkeys ${KEYMAP} windowkeys.inc
	eend $? "Error loading key mappings"
	return
}
