# init.d/keymaps
#!/sbin/runscript
# Copyright 1999-2002 Gentoo Technologies, Inc.
# Distributed under the terms of the GNU General Public License, v2 or later
# $Header$


depend() {
	need localmount
}

checkconfig() {
	if [ -z "$KEYMAP" ] ; then
		eerror "You need to setup KEYMAP first"
		return 1
	fi
}

start() {
	checkconfig || return 1
	ebegin "Loading key mappings"
	/usr/bin/loadkeys -q ${KEYMAP} windowkeys.inc
	eend $? "Error loading key mappings"
	return
}
