# init.d/local
#!/sbin/runscript
# Copyright 1999-2002 Gentoo Technologies, Inc.
# Distributed under the terms of the GNU General Public License, v2 or later
# $Header$


depend() {
	after *
}

start() {
	ebegin "Starting local"

	# Add any misc programs that should be started
	# to /etc/conf.d/local.start
	if [ -e /etc/conf.d/local.start ]
	then
		source /etc/conf.d/local.start
	fi
	
	eend $? "Started local."
}

stop() {
	ebegin "Stopping local"

	# Add any misc programs that should be stopped
	# to /etc/conf.d/local.stop

	if [ -e /etc/conf.d/local.stop ]
	then
		source /etc/conf.d/local.stop
	fi
	
	eend $? "Stopped local."
}
