# init.d/localmount
#!/sbin/runscript
# Copyright 1999-2002 Gentoo Technologies, Inc.
# Distributed under the terms of the GNU General Public License, v2 or later
# $Header$


depend() {
	need checkfs
}

start() {
	# LVM support for /usr, /home, /opt ....
	if [ -x /sbin/vgchange -a -f /etc/lvmtab ]
	then
		ebegin "Setting up the Logical Volume Manager"
		#still echo stderr for debugging
		/sbin/vgscan && /sbin/vgchange -a y >/dev/null
		eend $? "Failed to setup the LVM"
	fi
	
	# Mount local filesystems in /etc/fstab.
	ebegin "Mounting local filesystems"
	mount -at nocoda,nonfs,noproc,noncpfs,nosmbfs,noshm >/dev/null
	eend $? "Some local filesystem failed to mount"

	# Mount usbdevfs if kernel have support
	if [ -z "$(grep usbdevfs /proc/filesystems)" ]
	then
		# >/dev/null to hide errors from non-USB users
		modprobe usbcore &>/dev/null
	fi
	if [ -n "$(grep usbdevfs /proc/filesystems)" ] && [ -e /proc/bus/usb ] && \
	   [ ! -e /proc/bus/usb/devices ]
	then
		ebegin "Mounting USB device filesystem"
		mount -t usbdevfs usbdevfs /proc/bus/usb &>/dev/null
		eend $? "Failed to mount USB device filesystem"
	fi

	#swap on loopback devices, and other weirdnesses
	ebegin "Activating (possibly) more swap"
	/sbin/swapon -a &>/dev/null
	eend 0
}


# vim:ts=4
