# init.d/localmount
#!/sbin/runscript
# Copyright 1999-2004 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header$

depend() {
	need checkfs
}

start() {
	# Mount local filesystems in /etc/fstab.
	ebegin "Mounting local filesystems"
	mount -at nocoda,nonfs,noproc,noncpfs,nosmbfs,noshm,nocifs >/dev/null
	eend $? "Some local filesystem failed to mount"

	# Make sure we insert usbcore if its a module
	if [ -f /proc/modules ]
	then
		# >/dev/null to hide errors from non-USB users
		modprobe usbcore &>/dev/null
	fi
	
	# Check what USB fs the kernel support.  Currently
	# 2.5+ kernels, and later 2.4 kernels have 'usbfs',
	# while older kernels have 'usbdevfs'.
	local usbfs="$(grep -Fow usbfs /proc/filesystems || 
		grep -Fow usbdevfs /proc/filesystems)"

	if [ -n "${usbfs}" ] && \
	   [ -e /proc/bus/usb -a ! -e /proc/bus/usb/devices ]
	then
		ebegin "Mounting USB device filesystem (${usbfs})"
#		# Fetch usb gid from /etc/group; fixes bug 35860
#		usbgid=$(awk -F: '/^usb:/{print $3; exit}' /etc/group)
#		mount -t ${usbfs} ${usbgid:+-o devmode=0664,devgid=$usbgid} \
		mount -t ${usbfs} none /proc/bus/usb &>/dev/null
		eend $? "Failed to mount USB device filesystem"
	fi

	# Swap on loopback devices, and other weirdnesses
	ebegin "Activating (possibly) more swap"
	/sbin/swapon -a &>/dev/null
	eend 0

	# Run any post_mount commands for cryptfs

	if [ -f /etc/conf.d/cryptfs ]
	then
		ebegin "Running post_mount commands for cryptfs"

		/bin/egrep "^mount" /etc/conf.d/cryptfs | \
		while read mountline
		do
			mount=
			mount_point=
			post_mount=

			eval ${mountline}
 
			target=${mount}

			! /bin/cryptsetup status ${target}|egrep '\<active:' > /dev/null
			configured=$?

			if [ ${configured} -eq 1 ]
			then
				mount_point=`/bin/awk "/\/dev\/mapper\/${target}/ { print \\$2 }" /proc/mounts`
				if [ -n "${mount_point}" ]
				then
					if [ -n "${post_mount}" ]
					then
						if ! eval "${post_mount}" > /dev/null
						then
							ewarn "Failed to run post_mount commands on: ${target}"
						fi
					fi
				else
					ewarn "Failed to find mount point to ${target}. Skipping"
				fi
			else
				ewarn "Target ${target} wasn't mapped, skipping"
			fi

		done
	fi
}


# vim:ts=4
