# init.d/net
#!/sbin/runscript

#it'd be best to incorporate this functionality into 
#runscript.sh and have skeleton functions in here.

start() {
	ebegin "Starting network"
	eend
	local x
	for x in /etc/init.d/net.*
	do
		if [ -L /dev/shm/.init.d/${x##*/} ]
		then
			continue
		fi
		${x} start
	done
	return 0
}

stop() {
	ebegin "Stopping network"
	eend
	local x
	for x in /dev/shm/.init.d/net.*
	do
		if [ ! -L /dev/shm/.init.d/${x##*/} ]
		then
			continue
		fi
		${x} stop
	done
	return 0
}
