# init.d/net.lo
#!/sbin/runscript

IFACE=${1##*.}

start() {
	ebegin "Bringing ${IFACE} up"
	/sbin/ifconfig lo 127.0.0.1 up
	eend
}

stop() {
	ebegin "Bringing ${IFACE} down"
	/sbin/ifconfig $IFACE down
	eend
}
