# init.d/net.lo
#!/sbin/runscript
# Copyright 1999-2002 Gentoo Technologies, Inc.
# Distributed under the terms of the GNU General Public License v2
# $Header$


start() {
	ebegin "Bringing ${IFACE} up"
	/sbin/ifconfig lo 127.0.0.1 up &>/dev/null
	eend 0
}

stop() {
	ebegin "Bringing ${IFACE} down"
	/sbin/ifconfig ${IFACE} down &>/dev/null
	eend 0
}

