# init.d/net.ppp0
#!/sbin/runscript

PEER=MyPeer		# Define peer (aka ISP)
DEBUG=NO		# Turn on debugging
REDIAL=NO		# Redial after being dropped

start() {
	local CMD_LINE
	ebegin
	if [ "$DEBUG"="YES" ]; then
		CMD_LINE="-d"
	fi

	if [ "$REDIAL"="YES" ]; then
		CMD_LINE="$CMD_LINE -r"
	fi
	if [ -x "`which pon`" ]; then
		pon $CMD_LINE $PEER
	fi
	eend
}

stop() {
	ebegin
	if [ -x "`which poff`" ]; then
		/usr/sbin/poff $PEER
	fi
	eend
}
