# init.d/netmount
#!/sbin/runscript
# Copyright 1999-2006 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2

depend() {
	local myneed="net"
	local myuse=""

	# Only have Portmap as a dependency if there is a nfs mount in fstab 
	# that should be mounted at boot time.  Also filter out comments.
	local nfsmounts=$(awk '!/^#/ && ($3=="nfs" || $3=="nfs4") && $4 !~ /noauto/ { print $0 }' /etc/fstab)

	if [[ -n ${nfsmounts} ]] ; then
		if [[ -x /etc/init.d/rpcbind ]] ; then
			myneed="${myneed} rpcbind"
		else
			myneed="${myneed} portmap"
		fi
		myuse="${myuse} nfs nfsmount"
	else
		if [[ -x /etc/init.d/rpcbind ]] ; then
			myuse="${myuse} rpcbind"
		else
			myuse="${myuse} portmap"
		fi
	fi

	need ${myneed}
	use ${myuse}
}

remove_net_fs() {
	local fs
	rcfilesystems=" ${rcfilesystems} "
	for fs in "$@" ; do
		rcfilesystems=${rcfilesystems// ${fs} / }
	done
	rcfilesystems=${rcfilesystems# } # remove front and
	rcfilesystems=${rcfilesystems% } #   back spaces
}

start() {
	local rcfilesystems=${NET_FS_LIST}

	# Only try to mount NFS filesystems if portmap was started.
	# This is to fix "hang" problems for new users who do not
	# add portmap to the default runlevel.
	if ! service_started portmap && ! service_started rpcbind ; then
		remove_net_fs nfs nfs4
	fi
	rcfilesystems=${rcfilesystems// /,}   # convert to comma-separated

	ebegin $"Mounting network filesystems"
	mount -at ${rcfilesystems}

	if [[ $? != 0 ]] ; then
		ewend 1 "Could not mount all network filesystems!"
	else
		eend 0
	fi

	return 0
}

stop() {
	local rcfilesystems=${NET_FS_LIST}

	# We let the afs init script handle unmounting afs stuff
	# because it requires special handling of the afs daemon
	# and similar ugly cruft
	if service_started afs-client ; then
		remove_net_fs afs
	fi

	rcfilesystems=${rcfilesystems// /,}   # convert to comma-separated

	local ret
	ebegin $"Unmounting network filesystems"
	umount -at ${rcfilesystems}
	ret=$?
	eend ${ret} $"Failed to simply unmount filesystems"
	[[ ${ret} == 0 ]] && return 0

	# Nice function in rc-services.sh do handle the grunt work
	eindent
	do_unmount "umount" "" "" "^(${NET_FS_LIST// /|})$"
	eoutdent
	return ${retval}
}

# vim: set ts=4 :
