# init.d/netmount
#!/sbin/runscript
# Copyright 1999-2002 Gentoo Technologies, Inc.
# Distributed under the terms of the GNU General Public License, v2 or later
# $Header$


depend() {
	need net
	use portmap
}

start() {
	ebegin "Mounting network filesystems"
	mount -at coda,nfs,ncpfs,smbfs
	eend $?
}

stop() {
	ebegin "Unmounting network filesystems"
	# umount -art $fstypes doesn't seem to work, so...
	# NB: we have to check if any network filesystems is mounted, else mount complains

	local mymounts="`cat /proc/mounts | awk '{ print $3 " " $2 }' | grep -E ^'coda|nfs|ncpfs|smbfs' | awk 'BEGIN { ORS="	"} { print $2 }'`"

	if [ -n "$mymounts" ] ; then
		umount ${mymounts}
	fi
	eend $?
}
