# init.d/numlock
#!/sbin/runscript
# Copyright 1999-2002 Gentoo Technologies, Inc.
# Distributed under the terms of the GNU General Public License, v2 or later
# $Header$


depend() {
	need localmount
}

start() {
	ebegin "Enabling numlock on ttys"
	for tty in /dev/tty[1-8]
	do
		setleds -D +num < ${tty} &>/dev/null
	done
	eend $? "Failed to enable numlock"
}

stop() {
	ebegin "Disabling numlock on ttys"
	for tty in /dev/tty[1-8]
	do
		setleds -D -num < ${tty} &>/dev/null
	done
	eend $? "Failed to disable numlock"
}


# vim:ts=4
