# init.d/rmnologin
#!/sbin/runscript

depend() {
	need localmount
}

start() {
	if [ -f /etc/nologin.boot ]
	then
		rm -f /etc/nologin /etc/nologin.boot
	fi
	return
}
