# init.d/samba
#!/sbin/runscript
#RCUPDATE:3 4:76:This line is required for script management

. /etc/init.d/functions.sh


EXE0=/usr/sbin/nmbd
EXE1=/usr/sbin/smbd

depend() {
	need net
}

start() {
	ebegin "Starting samba"
	start-stop-daemon --start --quiet --exec $EXE0 
	result=$?
	start-stop-daemon --start --quiet --exec $EXE1 
	result=$(( $result + $? ))
	eend $result
}

stop() {
	ebegin "Stopping samba"
	start-stop-daemon --stop --quiet --exec $EXE0 
	result=$?
	start-stop-daemon --stop --quiet --exec $EXE1 
	result=$(( $result + $? ))
	eend $result
}




