# init.d/svc-xinetd
#!/sbin/runscript

SERVICE=xinetd

prepconfig() {
	if [ ! -e /etc/xinetd.conf ] ; then
		einfo "Creating new config from inetd..."
	  	/usr/sbin/xconv.pl < /etc/inetd.conf > /etc/xinetd.conf
	fi
	if [ /etc/inetd.conf -nt /etc/xinetd.conf ] ; then
		einfo "Creating new config from inetd..."
		/usr/sbin/xconv.pl < /etc/inetd.conf > /etc/xinetd.conf
	fi
	if [ ! -e /etc/xinetd.conf ]
	then
		eerror "xinetd not started.  Config file not found."
		exit 1
	fi
}
    
depend() {
	need net
}

start() {
	prepconfig
    ebegin "Starting supervised ${SERVICE}"
	ln -sf ../services/${SERVICE} ${SVCDIR}/control/${SERVICE}
    eend $?
}

stop() {
    ebegin "Stopping supervised ${SERVICE}"
   	if [ -e ${SVCDIR}/control/${SERVICE} ]
	then
		/usr/bin/svc -dx ${SVCDIR}/control/${SERVICE}
		rm ${SVCDIR}/control/${SERVICE}
	fi
	eend $?
}
