# init.d/svc-syslog
#!/sbin/runscript

SERVICE=syslog

start() {
    ebegin "Starting supervised ${SERVICE}"
	ln -sf ../services/${SERVICE} ${SVCDIR}/control/${SERVICE}
	eend $?
}

stop() {
    ebegin "Stopping supervised ${SERVICE}"
   	if [ -e ${SVCDIR}/control/${SERVICE} ]
	then
		/usr/bin/svc -dx ${SVCDIR}/control/${SERVICE}
		rm ${SVCDIR}/control/${SERVICE}
	fi
	eend $?
}
