# init.d/urandom
#!/sbin/runscript
# Copyright 1999-2002 Gentoo Technologies, Inc.
# Distributed under the terms of the GNU General Public License, v2 or later
# $Header$


depend() {
	need localmount
}

start() {
	[ -c /dev/urandom ] || return
	ebegin "Initializing random number generator"
	if [ -f /var/run/random-seed ]
	then
	    cat /var/run/random-seed >/dev/urandom >/dev/null 2>&1
	fi
	rm -f /var/run/random-seed
	umask 077
	dd if=/dev/urandom of=/var/run/random-seed count=1 >/dev/null 2>&1
	eend $? "Error initializing random number generator"
	umask 022
}

stop() {
	ebegin "Saving random seed"
	# Carry a random seed from shut-down to start-up;
	# see documentation in linux/drivers/char/random.c
	umask 077
	dd if=/dev/urandom of=/var/run/random-seed count=1 >/dev/null 2>&1
	eend $? "Could not save random seed."
}

